# supervisord.conf
[supervisord]
nodaemon = true

[program:avahi-daemon]
command=avahi-daemon

[program:dbus-daemon]
command=dbus-daemon --nofork --system

[program:homebridge]
command=npm start